# crontab
# “At 02:00 on every day-of-month”
0 2 * * * certbot renew --dry-run
